# traffic-masking/systemd/traffic-masking-server.service
[Unit]
Description=Traffic Masking Server
Documentation=https://github.com/kogeler/tooling/tree/main/traffic-masking
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=traffic-masking
Group=traffic-masking

WorkingDirectory=/opt/traffic-masking
Environment="PYTHONUNBUFFERED=1"
LoadCredential=control.psk:/etc/traffic-masking/control.psk

# Experimental native profile with an explicit ceiling
ExecStart=/usr/bin/python3 /opt/traffic-masking/traffic_masking_server.py \
    --host 0.0.0.0 \
    --port 8888 \
    --max-mbps 10 \
    --shape-mode profile \
    --profile mixed \
    --padding random \
    --mtu 1200 \
    --stats-interval 10 \
    --psk-file %d/control.psk \
    --max-clients 8 \
    --max-total-mbps 80 \
    --max-handshakes-per-second 20

Restart=on-failure
RestartSec=10
TimeoutStopSec=5

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077
LimitNOFILE=4096

StandardOutput=journal
StandardError=journal
SyslogIdentifier=traffic-masking-server

[Install]
WantedBy=multi-user.target
