# traffic-masking/systemd/traffic-masking-server.service
[Unit]
Description=Traffic Masking Server (Maximum Configuration)
Documentation=https://github.com/traffic-masking/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nobody
Group=nogroup

# Working directory
WorkingDirectory=/opt/traffic-masking

# Environment
Environment="PYTHONUNBUFFERED=1"

# Experimental native profile with an explicit ceiling
ExecStart=/usr/bin/python3 /opt/traffic-masking/traffic_masking_server.py \
    --host 0.0.0.0 \
    --port 8888 \
    --max-mbps 10 \
    --shape-mode profile \
    --profile mixed \
    --padding random \
    --mtu 1200 \
    --stats-interval 10 \
    --psk-file /etc/traffic-masking/control.psk \
    --max-clients 8 \
    --max-total-mbps 80 \
    --max-handshakes-per-second 20

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

# Resource limits
LimitNOFILE=65536
Nice=-10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=traffic-masking-server

[Install]
WantedBy=multi-user.target
